# distrib/dist/etc/init.d/sdwatchdog
#!/bin/sh

### BEGIN INIT INFO
# Provides:          sdwatchdog
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: speedata Publisher watchdog
# Description:       Starts the speedata Publisher watchdog to run the publishing process.
### END INIT INFO


set -e


DAEMON=/usr/bin/sdwatchdog
SDWATCHDOG_ENABLE=false
SDWATCHDOG_DEFAULTS_FILE=/etc/default/sdwatchdog


test -x $DAEMON || exit 0

. /lib/lsb/init-functions
. /etc/default/rcS

if [ -s $SDWATCHDOG_DEFAULTS_FILE ]; then
  . $SDWATCHDOG_DEFAULTS_FILE
  case "x$SDWATCHDOG_ENABLE" in
    xtrue|xfalse)   ;;
    *)    log_failure_msg "Value of SDWATCHDOG_ENABLE in $SDWATCHDOG_DEFAULTS_FILE must be either 'true' or 'false';"
          log_failure_msg "not starting sdwatchdog daemon."
          exit 1
          ;;
  esac
fi


case "$1" in
  start)
    if "$SDWATCHDOG_ENABLE"; then
      log_daemon_msg "Starting sdwatchdog daemon" "sdwatchdog"
      pkill -f 'texlua.*speedata.*watchdog' || :
      start-stop-daemon --start --exec $DAEMON --pidfile /var/run/watchdog.pid  -m -b -c speedata:speedata
    else
      if [ -s "$SDWATCHDOG_DEFAULTS_FILE" ]; then
          [ "$VERBOSE" != no ] && log_warning_msg "sdwatchdog daemon not enabled in /etc/default/sdwatchdog, not starting..."
      fi
    fi
    ;;
  stop)
    log_daemon_msg "Stopping sdwatchdog daemon" "sdwatchdog"
    start-stop-daemon --stop --pidfile /var/run/watchdog.pid -o
    rm -f /var/run/watchdog.pid
    ;;
  reload|force-reload|restart)
    if $SDWATCHDOG_ENABLE; then
      log_daemon_msg "Restarting sdwatchdog daemon" "sdwatchdog"
      start-stop-daemon --stop --pidfile /var/run/watchdog.pid  -o
      start-stop-daemon --start --exec $DAEMON --pidfile /var/run/watchdog.pid  -m -b -c speedata:speedata
    fi
    ;;
esac
exit 0
